# etc/supervisord.d/chronix.ini
[program:chronix]
command=
redirect_stderr=true
#stdout_logfile=syslog
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
priority = 10
